# scripts/jira-rag-sync.service
[Unit]
Description=Jira RAG incremental sync (cron safety-net for webhook)

[Service]
Type=oneshot
WorkingDirectory=/opt/jira-rag
EnvironmentFile=/opt/jira-rag/.env
ExecStart=/opt/jira-rag/scripts/run_sync.sh
